# technical_architecture/deployment/python-cognitive-agent.service
[Unit]
Description=Cognitive Intelligence Agent (Python Logic Layer)
Documentation=https://github.com/your-org/animal-vocalization
# Requires Rust field engine to be running first
Requires=rust-field-engine.service
After=rust-field-engine.service

[Service]
Type=simple
# Path to Python and the cognitive agent main script
ExecStart=/usr/bin/python3 /opt/cognitive/main.py
# Working directory for the cognitive agent
WorkingDirectory=/opt/cognitive
# Restart automatically on crash (Let it crash philosophy)
Restart=always
# Wait 2 seconds before restart (gives Rust time to detect disconnect)
RestartSec=2
# Standard output to journal for logging
StandardOutput=journal
StandardError=journal
# Security hardening
NoNewPrivileges=true
PrivateTmp=true
# Environment variables
Environment="PYTHONUNBUFFERED=1"
Environment="RUST_HEARTBEAT_ENDPOINT=ipc:///tmp/cognitive_heartbeat.ipc"

[Install]
WantedBy=multi-user.target
